# server/cs6238.conf
description "uWSGI server instance configured to server project 2 of CS6238"

start on runlevel [2345]
stop on runlevel [!2345]

setuid flask
setgid www-data

env PATH=/home/flask/cs6238-project2/proj2/bin
chdir /home/flask/cs6238-project2/server
exec uwsgi --ini server.ini
